# klatrebot.service
[Unit]
Description=KlatreBot Discord Bot Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_USER@
WorkingDirectory=@PROJECT_DIR@
EnvironmentFile=/etc/klatrebot/klatrebot.env
# Deploy step rewrites @POETRY_DIR@ with `dirname $(command -v poetry)` so systemd can find poetry.
Environment=PATH=@POETRY_DIR@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env poetry run python3 -m klatrebot_v2
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=klatrebot

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=@PROJECT_DIR@ @DATA_DIR@ /home/@SERVICE_USER@/.cache/pypoetry /home/@SERVICE_USER@/.config/pypoetry

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
