# klatrebot_v1/klatrebot.service
[Unit]
Description=KlatreBot Discord Bot Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=Admin
Group=Admin
WorkingDirectory=/home/Admin/KlatreBot/KlatreBot_Public
Environment=PATH=/home/Admin/KlatreBot/KlatreBot_Public/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/Admin/KlatreBot/KlatreBot_Public/.venv/lib/python3.11/site-packages
ExecStart=/usr/bin/python3 /home/Admin/KlatreBot/KlatreBot_Public/KlatreBot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=klatrebot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/Admin/KlatreBot/KlatreBot_Public

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
